# roles/example-web-service/templates/example_web_service_backend.service.j2
[Unit]
Description=example web service backend

[Service]
Type=simple
Restart=always
ExecStart={{ ansible_user_dir }}/example_web_service_backend.sh

[Install]
WantedBy=default.target
